# sshd.conf
[program:sshd]
directory=/usr/local/
command=/usr/sbin/sshd -D
stdout_logfile=/var/log/sshd.log
stderr_logfile=/var/log/sshd_error.log
autostart=true
autorestart=true
redirect_err=true
priority=1
